# provisioning/roles/devapp/templates/supervisor.conf
[program:{{ app_name }}]
directory={{ app_dir }}
command={{ venv_dir }}bin/python {{ app_dir }}manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
